# config/supervisor/ducishr.conf
[program:ducis_api]
user=admin
directory=/opt/apps/ducis_api/
#command=/root/.virtualenvs/susthome/bin/gunicorn susthome.wsgi -w 1 -g 1000 -u 1000 -t 300 -b 127.0.0.1:8002
command=/home/admin/.virtualenvs/ducis_api/bin/python manage.py runserver 127.0.0.1:8000
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=syslog
stderr_logfile=syslog
